# deploy/systemd/media-worker.service
[Unit]
Description=Translator Media Worker
After=network-online.target
Wants=network-online.target
PartOf=translator.target translator.service

[Service]
Type=simple
WorkingDirectory=/root/github/translator
EnvironmentFile=-/root/github/translator/.env
Environment=MEDIA_WORKER_HOST=127.0.0.1
Environment=MEDIA_WORKER_PORT=8090
Environment=NODE_ENV=production
ExecStart=/usr/bin/node services/media-worker/server.js
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
UMask=0077
ProtectSystem=full
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true

[Install]
WantedBy=translator.target
