# deploy/systemd/radio-worker.service
[Unit]
Description=Translator Live Radio Worker
After=network-online.target
Wants=network-online.target
PartOf=translator.target

[Service]
Type=simple
WorkingDirectory=/root/github/translator
Environment=RADIO_WORKER_HOST=0.0.0.0
Environment=RADIO_WORKER_PORT=8091
Environment=RADIO_STORAGE_DIR=/root/github/translator/data/radio
Environment=NODE_ENV=production
ExecStart=/usr/bin/node services/radio-worker/server.js
Restart=always
RestartSec=5

[Install]
WantedBy=translator.target
